# profile/template/includedir/etc/profile.d/editor.sh
export EDITOR=vi
#export EDITOR=vim
#export EDITOR=nvim
#export EDITOR=nano
#export EDITOR=micro
